# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:mailcatcher]
command=mailcatcher --foreground --ip 0.0.0.0

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
