# config/hrinfo/etc/supervisord.d/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/bin/bash -c "/usr/sbin/httpd -d /etc/httpd -f conf/httpd.conf -DFOREGROUND"
